# build/liberty-stack/systemd/liberty-stack-verify.timer
[Unit]
Description=Schedule Liberty Stack verification

[Timer]
OnCalendar=daily
Persistent=true
Unit=liberty-stack-verify.service

[Install]
WantedBy=timers.target
