# linux/systemd/sourceos-shell.service
[Unit]
Description=SourceOS Shell runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/sourceos-shell/bin/sourceos-shell
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
